# python_processing/background_worker.service
[Unit]
Description=Drone Crop Health Background Image Processing Worker
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Capstone_Interface/python_processing
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/home/ubuntu/Capstone_Interface/python_processing/.env
ExecStart=/usr/bin/python3 /home/ubuntu/Capstone_Interface/python_processing/background_worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

